# systemd/ssh-keep-awake.service
[Unit]
Description="Keep awake while SSH connection"

[Service]
ExecStart=/usr/local/bin/ssh-keep-awake.sh
